# scripts/services/m1-da-light-node/run
#!/usr/bin/env bash
set -e

set -a  # Automatically export all variables
source $DOT_MOVEMENT/env/celestia/.env
set +a  # Disable automatic export

./target/$CARGO_PROFILE/m1-da-light-node