# scripts/services/m1-da-light-node/run
#!/usr/bin/env bash
set -e

source $DOT_MOVEMENT/env/celestia/.env

./target/$CARGO_PROFILE/m1-da-light-node