# scripts/services/wait-for-celestia-light-node/run
#!/usr/bin/env bash
set -e

source $DOT_MOVEMENT/env/celestia/.env

./target/$CARGO_PROFILE/wait-for-celestia-light-node